# ansible/files/supervisord.conf
[program:imgresize]
command={{project_dir}}/imgresize -port 80
directory={{project_dir}}
user=www-data
autostart=true
autorestart=true
redirect_stderr=True
